# dot-config/systemd/user/verbum-terminal@.service
[Service]
Type=simple
ExecStart=/usr/libexec/gnome-terminal-server --class %I --app-id org.verbum.terminal.%I

[Install]
WantedBy=default.target
