# dot-config/systemd/user/verbum-sshagent.service
[Service]
Type=simple
ExecStart=/usr/bin/ssh-agent -D -a /run/user/1000/keyring/ssh

[Install]
WantedBy=default.target
